# files/inittab
::sysinit:/sbin/initscript
::shutdown:/bin/umount -a -r &> /dev/null
